# mailserver/target/filebeat.conf
description "Filebeat"
start on startup
stop on shutdown
respawn
respawn limit 20 5
task
exec /opt/filebeat/filebeat-linux-arm -c /etc/filebeat/filebeat.yml